# etc/supervisor_tensorboard.conf
[program:tensorboard]
command=tensorboard --logdir=/root/tmp/tensorboard
