# etc/supervisor_tensorboard.conf
[program:tensorboard]
command=tensorboard --logdir=src:/root/tmp/tensorboard/src,mnist:/root/tmp/tensorboard/mnist
