# etc/supervisor_gunicorn.conf
[program:webserver]
command=gunicorn --config=/etc/gunicorn.conf.py main:app
directory=/root/src/server
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
